# signalk-mqtt-bridge.service
[Unit]
Description=SignalK to MQTT Bridge
After=network.target signalk.service mosquitto.service
Wants=signalk.service mosquitto.service

[Service]
Type=simple
User=arielle
WorkingDirectory=/home/arielle/BoatOS
ExecStart=/usr/bin/python3 -u /home/arielle/BoatOS/signalk_to_mqtt.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
