# deploy/boatos-remote.service
[Unit]
Description=BoatOS Remote Control Server
After=network.target

[Service]
Type=simple
User=boatos
Group=input
WorkingDirectory=/home/boatos/BoatOS/backend
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/usr/bin/python3 /home/boatos/BoatOS/backend/remote_control_server.py
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boatos-remote

# Security
NoNewPrivileges=true
#PrivateTmp=true  # Disabled - causes issues with scrot screenshot capture

[Install]
WantedBy=multi-user.target
