# scripts/wifi-fallback.service
[Unit]
Description=BoatOS WiFi Fallback Hotspot
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
ExecStart=/bin/bash /home/boatos/BoatOS/scripts/wifi_fallback.sh
Restart=on-failure
RestartSec=30
User=root

[Install]
WantedBy=multi-user.target
